# emitter/snapper/snapper_supervisor.conf
[program:snapper]
command=/home/pi/whe/emitter/snapper/camcv_vid3
directory=/home/pi/whe/emitter
autostart=true
autorestart=true
redirect_stderr=true
logfile_maxbytes=1MB
logfile_backups=1