# printer/printer_supervisor.conf
# Prints to a local printer
[program:printer]
command=/opt/node/bin/node /home/pi/whe/printer/main.js
autostart=true
autorestart=true
redirect_stderr=true